# openrc/gpsmap
#!/sbin/openrc-run

supervisor=supervise-daemon

description="GPS Server"
command="/usr/local/bin/gpsmap"
command_args="-c /etc/gpsmap/gpsmap.toml"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/${RC_SVCNAME}.log"
error_log="/var/log/${RC_SVCNAME}.log"

depend() {
	need net
	after firewall
}